# bin/letterserver.service
[Unit]
Description=Mesh network signal service
Documentation=https://github.com/disktree/letterspace/blob/master/README.md
After=network.target

[Service]
WorkingDirectory=/usr/lib/letterserver
ExecStart=/usr/lib/letterserver/letterserver.js --service
Restart=always

[Install]
WantedBy=multi-user.target
